# infra/launchd/dk.broberg.trail.admin.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dk.broberg.trail.admin</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO__/scripts/trail-launchd-admin.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>RunAtLoad</key>
  <true/>

  <!-- Always restart on exit. To stop permanently: bash scripts/uninstall-launchd.sh -->
  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__HOME__/.trail/admin.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.trail/admin.log</string>

  <key>ProcessType</key>
  <string>Interactive</string>
</dict>
</plist>
